# systemd/sibsutis-schedule-bot.service
[Unit]
Description=Telegram bot for sibsutis-schedule
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Поправь пути под себя.
ExecStart=%h/.local/bin/sibsutis-schedule-bot --config %h/.config/sibsutis-schedule/config.txt
Restart=on-failure
RestartSec=10s
# Не уходить в перезапуски при ошибке конфига (exit 1).
SuccessExitStatus=0

[Install]
WantedBy=default.target
